# deploy/systemd/veridoc-error-monitor.timer
[Unit]
Description=Run VeriDoc runtime error monitor every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Unit=veridoc-error-monitor.service
Persistent=true

[Install]
WantedBy=timers.target
